# crontab
* * * * * /bin/bash /usr/src/river.rdf/indexer.sh | logger

# Mandatory blank line
