# crontab
* * * * * solo -port=6000 /bin/bash /usr/src/river.rdf/indexer.sh
# Mandatory blank line
